# service/assembly-archive-watcher.path.tmpl
[Path]
{{ASMARCHIVE_SYSTEMD_PATH_CHANGED_LINES}}

[Install]
WantedBy=default.target
